# src/DEBIAN/postinst
#!/bin/bash

launchctl load /System/Library/LaunchDaemons/com.chronic-dev.CDevReporter.plist


